# root/etc/services.d/healthcheck/run
#!/usr/bin/with-contenv bash

curl -Is $PING_URL
sleep 5m